# meta-pando-common/recipes-connectivity/openssh/pando-files/sshdgenkeys.service
[Unit]
Description=OpenSSH Key Generation
RequiresMountsFor=/var /run
Requires=pando-state.service bind-etc-ssh-hostkeys.service
After=pando-state.service bind-etc-ssh-hostkeys.service

[Service]
ExecStart=@LIBEXECDIR@/sshd_check_keys
Type=oneshot
RemainAfterExit=yes

[Install]                                                                                                                                                                                                                                                                                                                                                                                
WantedBy=multi-user.target
