# meta-pando-common/recipes-support/pando-config-vars/pando-config-vars/os-sshkeys.service
[Unit]
Description=OS SSH keys
Requires=pando-boot.service bind-home-root-.ssh.service
After=pando-boot.service bind-home-root-.ssh.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@SBINDIR@/os-sshkeys

[Install]
WantedBy=multi-user.target
