# meta-pando-common/recipes-core/systemd/systemd/pando.target
[Unit]
Description=Pando Target
Documentation=man:systemd.target(5)
